# .config/systemd/user/solaar.service
[Unit]
Description=Logitech devices managment

[Service]
ExecStart=/usr/bin/solaar --window=hide

[Install]
WantedBy=default.target
